# cron/crontab
* * * * * python3 /app/script.py >> /app/cron.log 2>&1
