# cron/crontab
*/5 * * * * python /app/save_pages.py > /app/cron.log 2>&1
